# vortex_server.service
[Unit]
Description=Vortex Media Server
After=network.target

[Service]
Type=simple
User=root
# Change the ExecStart path to wherever your compiled vortex-server binary is permanently located
WorkingDirectory=/opt/vortex
ExecStart=/usr/bin/vortex_server
# Optional: Set environment variables here if needed
Environment="RUST_LOG=info"
Environment="VORTEX_DATA_DIR=/opt/vortex/data"
Environment="VORTEX_STATIC_DIR=/opt/vortex/static"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
